# vendor/github.com/Cray-HPE/metal-provision/roles/ncn-common-setup/files/srv/cray/resources/metal/systemd/metal-iptables.service
[Unit]
Description=Loads Metal iptables config
After=local-fs.target network.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/iptables-restore /etc/iptables/metal.conf
Restart=no
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
